# raspberry/rizomtemp-backend.service
[Unit]
Description=Rizom Temp Backend
Documentation=https://github.com/seu-usuario/rizom-temp
After=network.target postgresql.service mosquitto.service
Requires=postgresql.service mosquitto.service

[Service]
Type=simple
User=rizomtemp
WorkingDirectory=/opt/rizomtemp/backend
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rizomtemp-backend

# Variáveis de ambiente — serão preenchidas pelo install.sh
EnvironmentFile=/opt/rizomtemp/.env

[Install]
WantedBy=multi-user.target
